# readonly_Documents/PowerShell/profile.d/10-initialization.ps1
(&mise activate pwsh --shims ) | Out-String | Invoke-Expression
Invoke-Expression (&starship init powershell)

Import-Module PSReadLine

Set-PSReadLineOption -EditMode Emacs