# readonly_Documents/PowerShell/profile.d/40-completion.ps1
kubectl completion powershell | Out-String | Invoke-Expression
Invoke-Expression -Command $(gh completion -s powershell | Out-String)
